# src/systemd_files/aarchup.timer
[Unit]
Description=Run foo every minute

[Timer]
OnBootSec=15m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
