# src/systemd_files/aarchup.service
[Unit]
Description=Run aarchup

[Service]
Type=simple
ExecStart=/usr/share/aarchup/noxaarchup.sh
