# src/systemd_files/aarchup.timer
[Unit]
Description=Run aarchup in a given frequency 

[Timer]
OnBootSec=15m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
